# templates/etc/profile.d/storm.sh.j2
export PATH=$PATH:{{storm_parent_dir}}/storm/bin
